# root-layer/etc/s6-overlay/s6-rc.d/svc-mod-wireguard-wghealth/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec s6-applyuidgid -u $(id -u abc) -g $(id -g abc) wghealth
